# systemd/clawos-toolbridge.service
[Unit]
Description=ClawOS toolbridge — tool execution boundary (filesystem, shell, API, browser)
Documentation=https://github.com/xbrxr03/clawos
After=clawos-policyd.service clawos-memd.service clawos-modeld.service
Requires=clawos-policyd.service clawos-memd.service clawos-modeld.service

[Service]
Type=simple
WorkingDirectory=%h/clawos
ExecStart=%h/clawos/venv/bin/python3 %h/clawos/services/toolbridge/main.py
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/clawos/logs/toolbridge.log
StandardError=append:%h/clawos/logs/toolbridge.log
Environment=PYTHONPATH=%h/clawos
Environment=CLAWOS_HOME=%h/clawos
Environment=CLAWOS_SERVICE=toolbridge

[Install]
WantedBy=default.target
