# systemd/clawos-agentd.service
[Unit]
Description=ClawOS agentd — task queue and agent session manager
Documentation=https://github.com/xbrxr03/clawos
After=clawos-policyd.service clawos-memd.service clawos-modeld.service clawos-toolbridge.service
Requires=clawos-policyd.service clawos-memd.service clawos-modeld.service clawos-toolbridge.service

[Service]
Type=simple
WorkingDirectory=%h/clawos
ExecStart=%h/clawos/venv/bin/python3 %h/clawos/services/agentd/service.py
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/clawos/logs/agentd.log
StandardError=append:%h/clawos/logs/agentd.log
Environment=PYTHONPATH=%h/clawos
Environment=CLAWOS_HOME=%h/clawos
Environment=CLAWOS_SERVICE=agentd

[Install]
WantedBy=default.target
