# systemd/clawos-dashd.service
[Unit]
Description=ClawOS dashd - canonical dashboard API and WebSocket event stream
Documentation=https://github.com/xbrxr03/clawos
After=clawos-clawd.service
Wants=clawos-clawd.service

[Service]
Type=simple
WorkingDirectory=%h/.clawos-runtime
ExecStart=%h/.clawos-runtime/venv/bin/python -m services.dashd.main
Restart=always
RestartSec=5
StandardOutput=append:%h/.clawos-runtime/logs/dashd.log
StandardError=append:%h/.clawos-runtime/logs/dashd.log
Environment=PYTHONPATH=%h/.clawos-runtime
Environment=CLAWOS_HOME=%h/.clawos-runtime
Environment=CLAWOS_SERVICE=dashd

[Install]
WantedBy=default.target
