# systemd/clawos-toolbridge.service
[Unit]
Description=ClawOS toolbridge — tool execution boundary (filesystem, shell, API, browser)
Documentation=https://github.com/xbrxr03/clawos
After=clawos-policyd.service clawos-memd.service clawos-modeld.service
Requires=clawos-policyd.service clawos-memd.service clawos-modeld.service

[Service]
Type=simple
WorkingDirectory=%h/.clawos-runtime
ExecStart=%h/.clawos-runtime/venv/bin/python3 %h/.clawos-runtime/services/toolbridge/main.py
Restart=always
RestartSec=3
StandardOutput=append:%h/.clawos-runtime/logs/toolbridge.log
StandardError=append:%h/.clawos-runtime/logs/toolbridge.log
Environment=PYTHONPATH=%h/.clawos-runtime
Environment=CLAWOS_HOME=%h/.clawos-runtime
Environment=CLAWOS_SERVICE=toolbridge

[Install]
WantedBy=default.target
