# systemd/clawos-gatewayd.service
[Unit]
Description=ClawOS gatewayd — OpenClaw WhatsApp bridge (port 18789)
Documentation=https://github.com/xbrxr03/clawos
After=clawos-clawd.service
Wants=clawos-clawd.service

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/usr/bin/openclaw gateway --allow-unconfigured
Restart=on-failure
RestartSec=10
# Give OpenClaw time to connect to WhatsApp
TimeoutStartSec=30
StandardOutput=append:%h/clawos/logs/gatewayd.log
StandardError=append:%h/clawos/logs/gatewayd.log
Environment=CLAWOS_HOME=%h/clawos
Environment=CLAWOS_SERVICE=gatewayd

[Install]
WantedBy=default.target
