# systemd/clawos-modeld.service
[Unit]
Description=ClawOS modeld — Ollama model routing and health checks
Documentation=https://github.com/xbrxr03/clawos
After=clawos-policyd.service
Requires=clawos-policyd.service

[Service]
Type=simple
WorkingDirectory=%h/.clawos-runtime
ExecStart=%h/.clawos-runtime/venv/bin/python3 %h/.clawos-runtime/services/modeld/main.py
Restart=always
RestartSec=5
StandardOutput=append:%h/.clawos-runtime/logs/modeld.log
StandardError=append:%h/.clawos-runtime/logs/modeld.log
Environment=PYTHONPATH=%h/.clawos-runtime
Environment=CLAWOS_HOME=%h/.clawos-runtime
Environment=CLAWOS_SERVICE=modeld
Environment=OLLAMA_HOST=http://localhost:11434

[Install]
WantedBy=default.target
