# systemd/clawos-policyd.service
[Unit]
Description=ClawOS policyd — permission gate, Merkle audit, lifecycle hooks
Documentation=https://github.com/xbrxr03/clawos
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.clawos-runtime
ExecStart=%h/.clawos-runtime/venv/bin/python3 %h/.clawos-runtime/services/policyd/main.py
Restart=always
RestartSec=3
StandardOutput=append:%h/.clawos-runtime/logs/policyd.log
StandardError=append:%h/.clawos-runtime/logs/policyd.log
Environment=PYTHONPATH=%h/.clawos-runtime
Environment=CLAWOS_HOME=%h/.clawos-runtime
Environment=CLAWOS_SERVICE=policyd

[Install]
WantedBy=default.target
