# systemd/clawos-a2ad.service
[Unit]
Description=ClawOS a2ad — Agent-to-Agent Protocol server
After=network.target clawos-agentd.service
Wants=clawos-agentd.service

[Service]
Type=simple
ExecStart=/usr/bin/env python3 -m services.a2ad.main
WorkingDirectory=%h/.clawos-runtime
Environment=PYTHONPATH=%h/.clawos-runtime
Restart=always
RestartSec=5

[Install]
WantedBy=clawos.target
