# packaging/systemd/clawos-clawd.service
[Unit]
Description=ClawOS clawd — Orchestration Daemon
After=clawos-agentd.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m services.clawd.service
WorkingDirectory=/opt/clawos
Restart=always
RestartSec=3
Environment=PYTHONPATH=/opt/clawos

[Install]
WantedBy=default.target
