# dashboard/backend/dashd@.service
[Unit]
Description=ClawOS Dashboard Service (dashd)
After=network.target agentd.service policyd.service memd.service modeld.service
Wants=agentd.service policyd.service memd.service modeld.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/clawos/dashboard/backend
ExecStart=%h/.local/bin/uvicorn service:app --host 0.0.0.0 --port 7070 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=%h/clawos

[Install]
WantedBy=default.target
